# layers/meta-testdistro/recipes-core/tegra-startup/tegra-startup/set-tegra-clocks.timer
[Unit]
Description=Tegra clock setting delay

[Timer]
OnBootSec=60

[Install]
WantedBy=timers.target
