# layers/meta-testdistro/recipes-core/tegra-startup/tegra-startup/set-tegra-clocks.service.in
[Unit]
Description=Set Tegra clocks
After=nvstartup.service
Before=xserver-nodm.service weston.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@SBINDIR@/set-tegra-clocks

[Install]
WantedBy=multi-user.target
