# layers/meta-testdistro/recipes-core/tegra-startup/tegra-startup/set-tegra-clocks.service.in
[Unit]
Description=Set Tegra clocks

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@SBINDIR@/set-tegra-clocks

[Install]
WantedBy=multi-user.target
